# apps/relay/deploy/mixdog-relay.service
# systemd unit for the Mixdog relay (installed by deploy.sh to
# /etc/systemd/system/mixdog-relay.service).
[Unit]
Description=Mixdog remote relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mixdog-relay
WorkingDirectory=/opt/mixdog-relay
Environment=PORT=443
Environment=DATA_DIR=/var/lib/mixdog-relay
Environment=RENDERER_DIR=/opt/mixdog-relay/renderer
Environment=TLS_CERT=/etc/letsencrypt/live/RELAY_DOMAIN/fullchain.pem
Environment=TLS_KEY=/etc/letsencrypt/live/RELAY_DOMAIN/privkey.pem
ExecStart=/usr/bin/node /opt/mixdog-relay/server.mjs
Restart=always
RestartSec=3
# Thousands of concurrent phones/desktops need more than the 1024 default
# socket descriptors; per-service traffic shows up in systemctl show.
LimitNOFILE=65535
IPAccounting=yes
# Runaway protection on a 1GB box: cap the node heap below the cgroup limit
# so pressure surfaces as GC/backpressure (or a clean restart), never the
# host OOM killer taking the whole machine down.
Environment=NODE_OPTIONS=--max-old-space-size=512
MemoryMax=800M
# Bind :443 without root.
AmbientCapabilities=CAP_NET_BIND_SERVICE
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/mixdog-relay
ProtectHome=true

[Install]
WantedBy=multi-user.target
